# config/s6/next-terminal/run
#!/bin/sh
cd /usr/local/next-terminal/ && exec /usr/local/next-terminal/next-terminal -c /etc/next-terminal/config.yaml
